# work/crontab
# min hour monthday month weekday command
* * * * * ipython
